# deployment/conf/production/production_gunicorn.service
[Unit]
Description=production-gunicorn
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
RuntimeDirectory=gunicorn
EnvironmentFile=/mnt/data/var/vhosts/production/.env
WorkingDirectory=/mnt/data/var/vhosts/production/backend/src
ExecStart=/mnt/data/var/vhosts/production/venv/bin/gunicorn wagtailnest.wsgi:application \
  --workers=2 --log-level error \
  --bind=unix:/var/run/gunicorn/gunicorn.socket

[Install]
WantedBy=multi-user.target
